# libreelec/rpi2/tuning_audio.service
[Unit]
Description=Tuning Audio
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/libreelec/tuning_audio
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
