# libreelec/rpi2/tuning_audio.target
[Unit]
Description=Tuning Audio
Requires=tuning_audio.service
After=tuning_audio.service
AllowIsolate=yes
